# templates/user-mounts/systemd.mount.j2
[Unit]
Description={{ mount.description }}

[Install]
WantedBy=default.target

[Mount]
What={{ mount.what }}
Where={{ mount.where }}

Type={{ mount.type }}
