# templates/periodic-upgrade/systemd.timer.j2
# {{ ansible_managed }}
[Unit]
Description=Upgrade and restart
[Timer]
OnCalendar={{ upgrade_schedule }}
Persistent=true
[Install]
WantedBy=timers.target
