# templates/systemd.timer.j2
# {{ ansible_managed }}
[Unit]
Description={{ timer.name }}
[Timer]
OnCalendar={{ timer.time }}
Persistent=true
[Install]
WantedBy=timers.target
